# scripts/cpu_agent.service.example
[Unit]
Description=LLM Orchestration CPU Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bryan
WorkingDirectory=/home/bryan/llm_orchestration
ExecStart=/usr/bin/python3 /media/bryan/shared/scripts/cpu_agent.py --config /media/bryan/shared/agents/config.json --name cpu-%H
Restart=always
RestartSec=2
StandardOutput=append:/media/bryan/shared/logs/cpu_workers/cpu-agent-%H.log
StandardError=append:/media/bryan/shared/logs/cpu_workers/cpu-agent-%H.log

[Install]
WantedBy=multi-user.target
